# river/var/sv/pipewire-pulse/run
#!/bin/sh

sv check ../pipewire || exit 1

exec pipewire-pulse
