# wayland/var/sv/xdg-desktop-portal/run
#!/bin/sh

exec /usr/lib/xdg-desktop-portal -r
