# wayland/var/sv/swaybg/run
#!/bin/sh

exec swaybg --mode fill --image $HOME/.config/background
