# river/var/sv/swaybg/run
#!/bin/sh

exec swaybg -i $HOME/.config/background
