# river/var/sv/xdg-desktop-portal-wlr/run
#!/bin/sh

exec /usr/lib/xdg-desktop-portal-wlr -r
